# appstore-override/cronic/crontab
# Run every minute
*/1 * * * * echo "hello"

# Run every 2 seconds
docker ps -a

# Run once every hour
@hourly echo "$SOME_HOURLY_JOB"